# 1c/scripts/99-aladdin.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="0529", ATTR{idProduct}=="0001", RUN+="/usr/local/bin/link-aladdin"